# systemd/user/spotdl-navidrome-sync.service
[Unit]
Description=Sync Spotify and YouTube Music sources into Navidrome
Documentation=file:/home/aym-n/homelab/README.md

[Service]
Type=oneshot
WorkingDirectory=/home/aym-n/homelab
Environment=SPOTDL_THREADS=3
Environment=SPOTDL_MAX_RETRIES=2
Environment=YTMUSIC_LIMIT=5000
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
MemoryHigh=1G
MemoryMax=1536M
ExecStart=/home/aym-n/homelab/scripts/spotdl-navidrome-sync.sh
