# systemd/user/spotdl-navidrome-sync.timer
[Unit]
Description=Hourly Navidrome music source sync

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h
Persistent=true
Unit=spotdl-navidrome-sync.service

[Install]
WantedBy=timers.target
